# src-tauri/packaging/pipe-deck-daemon.service
# pipe-deck-daemon-unit-version: 2
[Unit]
Description=Pipe Deck background restore daemon
After=pipewire.service wireplumber.service pipewire-pulse.service
Wants=pipewire.service wireplumber.service

[Service]
Type=notify
NotifyAccess=main
Restart=on-failure
ExecStart=/usr/bin/pipe-deck-daemon
Environment=PIPE_DECK_USE_MOCK=0

[Install]
WantedBy=default.target
